# etc/meme.linux.in
#!/bin/sh
#
# chkconfig: 345 90 90
# Startup script to be placed as /etc/init.d/meme
#
# Starts/stops Meme and Mast servers 
#
# This file is automatically generated by make 
# $Id: meme.linux.in,v 1.2 2005/10/07 19:18:50 nadya Exp $
# $Log: meme.linux.in,v $
# Revision 1.2  2005/10/07 19:18:50  nadya
# soruce meme_config to get needed environment
# check for a specific user to start servers
#
# Revision 1.1  2005/09/26 23:15:16  nadya
# add startup scripts for linux and solaris.
#
#

# Source function library.
. /etc/init.d/functions

# site-specific configuration file 
CONF=@MEMEDIR@/etc/meme_config

if [ -f $CONF ] ; then
    . $CONF
else
    echo "$CONF does not exist. Meme installation is incomplete"
    exit 1
fi

# check for user existence
checkUser() {
    USER=`id -un $MEME_USER 2> /dev/null`
    if [ -z "$USER" ] ; then
        echo "User $MEME_USER does not exist. Meme installation is incomplete"
        exit 1
    fi
}

case "$1" in
    start)
        checkUser
        if [ -x $MEME_SERVER -a -x $MAST_SERVER ]; then
            action $"Starting Meme server" su $MEME_USER -c "$MEME_SERVER" 2> /dev/null
            action $"Starting Mast server" su $MEME_USER -c "$MAST_SERVER" 2> /dev/null
            RETVAL=$?
        else
            echo "$MEME_SERVER and $MAST_SERVER are not executables."
            echo "Skipping Meme and Mast startup"
            RETVAL=1
        fi
        ;;
    stop)
         echo -n $"Shutting down Meme server" 
         killproc meme-server
         echo
         echo -n $"Shutting down Mast server" 
         killproc mast-server
         echo
         RETVAL=$?
        ;;
    restart)
        $0 stop
        $0 start
        RETVAL=$?
        ;;
    *)
        echo "Usage: $0 { start | stop | restart }"
        exit 1
esac

exit $RETVAL

